# bin/hornet-firewall.service
[Unit]
Description=Hornet agent firewall rules
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=__REPO_DIR__/bin/setup-firewall.sh
ExecStop=/bin/sh -c 'UID_HORNET=$(id -u hornet_agent) && iptables -w -D OUTPUT -m owner --uid-owner $UID_HORNET -j HORNET_OUTPUT 2>/dev/null; iptables -w -F HORNET_OUTPUT 2>/dev/null; iptables -w -X HORNET_OUTPUT 2>/dev/null'

[Install]
WantedBy=multi-user.target
